# root/etc/s6-overlay/s6-rc.d/svc-code-server/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

mkdir -p /config/logs
exec 1> >(tee -a /config/logs/code-server.log)
exec 2> >(tee -a /config/logs/code-server.log >&2)

echo "[$(date)] Starting code-server service"

# Check if code-server binary exists
if [[ ! -f /app/code-server/bin/code-server ]]; then
    echo "[$(date)] ERROR: code-server binary not found at /app/code-server/bin/code-server"
    exit 1
fi

# Check if config directories exist and are writable
for dir in "/config/data" "/config/extensions" "/workspace"; do
    if [[ ! -d "$dir" ]]; then
        echo "[$(date)] ERROR: Directory $dir does not exist"
        exit 1
    fi
    if [[ ! -w "$dir" ]]; then
        echo "[$(date)] ERROR: Directory $dir is not writable"
        ls -la "$dir"
        exit 1
    fi
done

if [[ -n "${PASSWORD}" ]] || [[ -n "${HASHED_PASSWORD}" ]]; then
    AUTH="password"
    echo "[$(date)] Authentication mode: password"
else
    AUTH="none"
    echo "[$(date)] Starting with no password"
fi

if [[ -z ${PROXY_DOMAIN+x} ]]; then
    PROXY_DOMAIN_ARG=""
    echo "[$(date)] No proxy domain configured"
else
    PROXY_DOMAIN_ARG="--proxy-domain=${PROXY_DOMAIN}"
    echo "[$(date)] Using proxy domain: ${PROXY_DOMAIN}"
fi

echo "[$(date)] System Information:"
echo "User running process: $(whoami)"
echo "Current working directory: $(pwd)"
echo "Directory permissions:"
ls -la /config/
ls -la /workspace/

echo "[$(date)] Starting code-server with following configuration:"
echo "- Bind address: 0.0.0.0:8443"
echo "- User data dir: /config/data"
echo "- Extensions dir: /config/extensions"
echo "- Auth mode: ${AUTH}"
echo "- Proxy domain arg: ${PROXY_DOMAIN_ARG}"
echo "- Default workspace: ${DEFAULT_WORKSPACE:-/config/workspace}"

# Test if port is already in use
if nc -z 127.0.0.1 8443 2>/dev/null; then
    echo "[$(date)] ERROR: Port 8443 is already in use"
    netstat -tulpn | grep 8443
    exit 1
fi

exec \
    s6-notifyoncheck -d -n 300 -w 1000 -c "nc -z 127.0.0.1 8443" \
        s6-setuidgid abc \
            /app/code-server/bin/code-server \
                --bind-addr 0.0.0.0:8443 \
                --user-data-dir /config/data \
                --extensions-dir /config/extensions \
                --disable-telemetry \
                --auth "${AUTH}" \
                --verbose \
                "${PROXY_DOMAIN_ARG}" \
                "${DEFAULT_WORKSPACE:-/config/workspace}"
